# laptasus/laptasus.in
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="laptasus init script"
proj=%PROJ%
bindir=%BINDIR%
command=$bindir/$proj

depend() {
    need tmpfiles.setup
}

start() {
   ebegin "Starting $0"
   nohup $command start > /var/log/laptasus.log 2>&1 &
   eend $?
}

stop() {
    ebegin "Rejecting new connections."
    $command stop
    eend $? "Failed to stop daemon."
}
